# files/usr/lib/systemd/system/nix-daemon.service.d/mount.conf
[Unit]
Requires=nix.mount nix-store-init.service
After=nix.mount nix-store-init.service
