# files/usr/lib/systemd/system/nix-daemon.service.d/mount.conf
[Unit]
Requires=nix.mount
After=nix.mount
